# service/usbkill.service
[Unit]
After=network.service

[Service]
ExecStart=/root/.local/bin/usbkill

[Install]
WantedBy=default.target
